# vault/rootfs/etc/services.d/vault-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Vault
# ==============================================================================

bashio::log.info "Starting Vault Server..."

exec /usr/bin/vault-server.sh